# systemd/hmg-soar-report.service
[Unit]
Description=HMG Wazuh SOAR Brain - Geração de Relatório de Vulnerabilidades
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=hmg-soar
Group=www-data
WorkingDirectory=/opt/hmg-soar
EnvironmentFile=/etc/hmg-soar/credentials.env

# IMPORTANTE: --mode audit garante que NUNCA haverá correção automática
ExecStart=/usr/bin/python3 /opt/hmg-soar/analyserV1.py \
    --all-agents \
    --mode audit \
    --output /opt/hmg-soar/output/relatorio_wazuh.csv \
    --pdf /opt/hmg-soar/output/relatorio_wazuh.pdf \
    --html /opt/hmg-soar/output/relatorio_wazuh.html \
    --web-output-dir /var/www/wazuh-soar \
    --web-group www-data

StandardOutput=journal
StandardError=journal

# ====================================================================
# HARDENING (compatível: precisa de REDE para Wazuh/Indexer e escrita
# em output/cache/web). Sem MemoryDenyWriteExecute (libs nativas: pdf/csv).
# ====================================================================
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native
UMask=0027

# Sem capabilities especiais.
CapabilityBoundingSet=
AmbientCapabilities=

# Gravação restrita aos diretórios realmente necessários.
ReadWritePaths=/opt/hmg-soar/output /opt/hmg-soar/.hmg_cache /opt/hmg-soar/config /var/www/wazuh-soar /var/www/wazuh-soar/data /var/www/wazuh-soar/reports /var/www/wazuh-soar/assets

[Install]
WantedBy=multi-user.target
